# .claude/autonomous-refactor/chimera-autonomous-refactor.service
[Unit]
Description=Project Chimera Autonomous Codebase Refactoring Loop
Documentation=https://github.com/project-chimera/docs/autonomous-refactoring.md
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=ranj
Group=ranj
WorkingDirectory=/home/ranj/Project_Chimera
Environment="PATH=/home/ranj/.local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/ranj/Project_Chimera"
Environment="CHIMERA_BASELINE_FILE=/home/ranj/Project_Chimera/baseline_metrics.json"
Environment="CHIMERA_MIN_COVERAGE=80.0"
Environment="CHIMERA_MUTATION_TESTS=false"

# Resource limits to prevent interference with foreground development
MemoryMax=4G
CPUQuota=20%
TasksMax=100

# Restart policy
Restart=on-failure
RestartSec=60

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/ranj/Project_Chimera
ReadWritePaths=/home/ranj/Project_Chimera/.claude/autonomous-refactor
ReadWritePaths=/home/ranj/Project_Chimera/baseline_metrics.json
ReadWritePaths=/home/ranj/Project_Chimera/evaluation_results.json

# Logging
StandardOutput=append:/home/ranj/.claude/autonomous-refactor/loop.log
StandardError=append:/home/ranj/.claude/autonomous-refactor/loop_errors.log

# Main command - run the orchestration script
ExecStart=/home/ranj/.local/bin/python3 services/autonomous-agent/orchestrator.py

# Graceful shutdown
ExecStop=/bin/kill -SIGTERM $MAINPID
TimeoutStopSec=30
KillMode=mixed

[Install]
WantedBy=multi-user.target
